# systemd/llama-second-opinion.service
[Unit]
Description=llama-server for second-opinion (Qwen3-Coder-30B on 7900 XTX)
# No reverse dependency on codium: PartOf/BindsTo both caused llama to
# stop immediately at startup because codium hadn't yet been active in
# the transaction. Auto-teardown by systemd is deferred (see ha-fix
# branch); the launcher's wait loop stops llama explicitly after codium
# exits.

[Service]
Type=exec
ExecStart=%h/Documents/Repos/Workstation/second-opinion/scripts/primary-llama-vulkan.sh
KillSignal=SIGTERM
TimeoutStopSec=10
Restart=no
StandardOutput=journal
StandardError=journal
