# systemd/llama-primary-experiment.service
[Unit]
Description=llama-server EXPERIMENT (GPT-OSS-120B on 7900 XTX / HIP, port 11444, MoE expert offload)
# EXPERIMENTAL — MoE expert offload investigation, phase 2.
# Runs mainline ggml-org/llama.cpp built with HIP (ROCm 7.2.1) for gfx1100.
# Stock Vulkan path failed: RADV cannot allocate the 45 GB host-visible
# Vulkan buffer that --n-cpu-moe requires for expert weights. HIP uses
# hipHostMalloc and doesn't have the same per-allocation ceiling.
#
# (ik_llama.cpp fork was the original phase-2 plan but its HIP build is
# bit-rotted under modern ROCm — see docs/research/gpt-oss-120b-moe-offload.md
# for the full investigation log.)
#
# Binary at /usr/local/lib/llama.cpp-hip/llama-server, separate from the
# stock Vulkan /usr/local/lib/llama.cpp/ used by llama-primary.service.
#
# Cannot coexist with llama-primary.service (both want Vulkan0). Workflow:
#   systemctl stop llama-primary
#   systemctl start llama-primary-experiment
#   ...test, measure tok/s...
#   systemctl stop llama-primary-experiment
#   systemctl start llama-primary
#
# Tunables to sweep:
#   --n-cpu-moe N    start at 26 (the Register's value for a 20GB GPU)
#                    decrease until OOM, then back off. With 24GB on the 7900
#                    XTX we should be able to push N lower than 26 — meaning
#                    fewer expert layers exiled to CPU, more on the GPU,
#                    higher tok/s.
#   -c CONTEXT       start at 32768 to match the article; bump to 65536 or
#                    90000 once N is settled if KV cache fits.
#   -ub / -b         larger batches help CPU+GPU hybrid; default 4096 is a
#                    reasonable starting point for MoE hybrid.
#
# Sampler flags here are GPT-OSS-recommended starting points — re-derive if
# the model changes.
#   - temp 1.0, top_p 1.0, top_k 0, min_p 0.0 (GPT-OSS recommended)
#   - Adjust per Hugging Face model card if needed.

[Service]
Type=simple
EnvironmentFile=/etc/workstation/system.env
# HIP_VISIBLE_DEVICES=1 → only expose the 7900 XTX (gfx1100) to llama-server.
# CRITICAL: HIP enumerates differently from rocminfo. rocm-smi and HIP runtime
# both put the 5700 XT at index 0 and the 7900 XTX at index 1, while rocminfo
# puts them in the reverse order. Without this constraint, llama-server picks
# the 5700 XT (gfx1010) by default, which we did NOT compile kernels for.
Environment="HIP_VISIBLE_DEVICES=1"
Environment="LLAMA_SET_ROWS=1"
ExecStart=/usr/local/lib/llama.cpp-hip/llama-server \
  -m ${WS_MODELS_DIR}/gpt-oss-120b/gpt-oss-120b-mxfp4-00001-of-00003.gguf \
  --alias gpt-oss-120b \
  -ngl 999 --n-cpu-moe 28 \
  -c 131072 --flash-attn on \
  -ctk q8_0 -ctv q8_0 \
  --jinja \
  -b 2048 -ub 2048 \
  --cache-reuse 256 \
  --no-mmap \
  --temp 1.0 --top-p 1.0 --top-k 0 --min-p 0.0 \
  --host 127.0.0.1 --port 11444
Restart=no
KillSignal=SIGTERM
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal

# Restart=no on purpose: this is an experiment. If it OOMs or crashes we want
# to see the failure and adjust flags, not have systemd loop-restart it.
