# archive/library_and_patron/archive/v1-code-oriented/systemd/librarian.service
[Unit]
Description=Librarian MCP Server — RAG context server
After=network.target ollama-gpu0.service ollama-gpu1.service
Wants=ollama-gpu0.service ollama-gpu1.service
# Do NOT enable this service — managed by ai-session (session manager)

[Service]
User=levine
Group=levine
WorkingDirectory=/home/levine/Documents/Repos/Workstation/library_and_patron
ExecStart=/home/levine/Documents/Repos/Workstation/library_and_patron/.venv/bin/python -m librarian.server
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment="ROCR_VISIBLE_DEVICES=0"
Environment="HSA_OVERRIDE_GFX_VERSION=10.3.0"
Environment="OLLAMA_GPU1_URL=http://localhost:11434"
Environment="OLLAMA_GPU0_URL=http://localhost:11435"
Environment="LIBRARIAN_PORT=11436"
Environment="LIBRARIAN_DB_BASE=/home/levine/.local/share/librarian/lancedb"
Environment="LIBRARIAN_WATCH_CONFIG=/home/levine/.config/librarian/watched_repos.yaml"

[Install]
WantedBy=multi-user.target
